# pi-cecilia/blackroad-operator/blackroad-sf/blackroad-sf/service/myapp.service
[Unit]
Description=My App Service
After=network.target

[Service]
ExecStart=/usr/bin/env bash /opt/myapp/app.sh
Restart=always
EnvironmentFile=/opt/myapp/config.env

[Install]
WantedBy=multi-user.target
