# pi-lucidia/octavia-operator/packages/cece-cli/services/cece-api.service
[Unit]
Description=CECE API Gateway
Documentation=https://blackroad.io/cece
After=network.target

[Service]
Type=simple
User=cecilia
WorkingDirectory=/home/cecilia/cece
ExecStart=/usr/bin/node /home/cecilia/cece/api/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOME=/home/cecilia

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cece-api

[Install]
WantedBy=multi-user.target
