# pi-lucidia/octavia-operator/packages/cece-cli/services/cece-heartbeat.service
[Unit]
Description=CECE Heartbeat Monitor
Documentation=https://blackroad.io/cece
After=network.target

[Service]
Type=simple
User=cecilia
WorkingDirectory=/home/cecilia/cece/heart
ExecStart=/home/cecilia/cece/heart/heartbeat.sh
Restart=always
RestartSec=60
Environment=HOME=/home/cecilia

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cece-heartbeat

[Install]
WantedBy=multi-user.target
