# pi-lucidia/lucidia-home/lucidia.service
[Unit]
Description=Autonomous Lucidia - BlackRoad Guardian
After=network.target

[Service]
Type=simple
User=lucidia
WorkingDirectory=/home/lucidia
ExecStart=/usr/bin/python3 /home/lucidia/lucidia-autonomous.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
